# trmnl-terminus/rootfs/etc/s6-overlay/s6-rc.d/postgresql/finish
#!/usr/bin/bash

echo "[postgresql] Service exited with code $1"