# trmnl-terminus/rootfs/etc/s6-overlay/s6-rc.d/terminus/finish
#!/usr/bin/bash

echo "[terminus] Service exited with code $1"