# trmnl-terminus/rootfs/etc/s6-overlay/s6-rc.d/init-db/up
#!/command/execlineb -P
/usr/bin/bash /etc/s6-overlay/scripts/init-db.sh